# systemd/nbh-monitor.service
# this is meant to be installed under /etc/systemd/system
[Unit]
Description=check for jupyter instances that have been running for too long

# can use e.g. --grace=180 to increase grace time to 3 hours
# or --period=5 if we want to run every 5 minutes
[Service]
Environment=PYTHONPATH=/root/nbhosting/django
ExecStart=/bin/bash -c "/usr/bin/python3 /root/nbhosting/django/stats/monitor.py"

[Install]
WantedBy=multi-user.target
