# systemd/nbh-uwsgi.service
# this is meant to be installed under /etc/systemd/system
[Unit]
Description=the nbhosting django app as a uswgi instance

[Service]
ExecStart=/usr/sbin/uwsgi --ini /etc/uwsgi.d/nbhosting.ini

[Install]
WantedBy=multi-user.target
