# systemd/nbh-monitor.service.in
# this is meant to be installed under /etc/systemd/system
[Unit]
Description=job that cyclically collects system and course data, and kills idle containers

# for complete help and default values, run
# nbh-manage monitor --help
#
[Service]
ExecStart=/bin/bash -c "python3 /usr/bin/nbh-manage monitor --period @monitor_period@ --idle @monitor_idle@"

[Install]
WantedBy=multi-user.target
